# data/hp-omen-key.service
[Unit]
Description=HP Omen Key Handler — launches HP Manager on Omen Key press
After=multi-user.target hp-manager.service
Wants=hp-manager.service

[Service]
Type=simple
ExecStart=/usr/libexec/hp-manager/omen-key-listener.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
