# data/hpm-mux.service
[Unit]
Description=OMEN Command Center — MUX (GPU Switch) Service
After=dbus.service
Requires=dbus.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/libexec/hp-manager/services/mux_service.py
Restart=on-failure
RestartSec=3

# Security hardening — namespace-creating directives (PrivateTmp, ProtectHome)
# disabled to fix status=226/NAMESPACE on Ubuntu 24.04+/Zorin OS 18.
# These conflict with AppArmor profiles or systemd configurations on affected kernels.
# Non-namespace directives are kept for defence-in-depth.
ProtectHome=false
PrivateTmp=false
NoNewPrivileges=true
LockPersonality=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
